# systemd/nextdns-ddns.service
[Unit]
Description=NextDNS DDNS Updater

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-nextdns.sh
